# sv/zookeeper/run
#!/bin/bash

#sv start socklog-unix || exit 1

source /etc/envvars
env

cd /zookeeper/bin

exec 2>&1
exec ./zkServer.sh start-foreground
